# gunicorn.conf
[program:gunicorn]
command=/envs/data_hub/bin/gunicorn data_hub.wsgi -b 0.0.0.0:8000
directory=/src/data_hub
